# com.user.dotfiles.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.dotfiles</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/bhushan/Desktop/dev1/dotfiles/dotfiles.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/bhushan/Desktop/dev1/dotfiles/dotfiles.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/bhushan/Desktop/dev1/dotfiles</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/bhushan/Desktop/dev1/dotfiles/sync.sh</string>
    </array>
</dict>
</plist>

